# agent/install/llama_server.service.tmpl
[Unit]
Description=llama.cpp server (llama-server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
# Add model flags before enabling, e.g.:  -hf <repo>:<quant>   or   -m /path/model.gguf
# --metrics is required for the agent to scrape llama-server telemetry.
ExecStart=__LLAMA_BIN__ --metrics --host 127.0.0.1 --port 8080
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
